# kresus/rootfs/etc/services.d/kresus/run
#!/usr/bin/with-contenv bashio
WOOB_DIR="/woob"
KRESUS_DATA_DIR="/data/kresus_data"
KRESUS_USER="kresus"
KRESUS_INI_FILE="/etc/kresus/config.ini"
KRUN="su-exec ${KRESUS_USER}:${KRESUS_USER}"

# ==============================================================================
# Pull latest Woob version
# ==============================================================================
cd "${WOOB_DIR}" || bashio::exit.nok

if ! bashio::fs.directory_exists "${WOOB_DIR}/.git"; then
    bashio::log.info "Installing woob..."
    ${KRUN} git clone --depth 1 https://gitlab.com/woob/woob.git . || (bashio::log.error "Couldn't install woob" && bashio::exit.nok)
    bashio::log.info "Done installing woob."
else
    bashio::log.info "Updating woob..."
    ${KRUN} git pull || bashio::log.warning "Couldn't update; maybe the Woob's server is unreachable?" && bashio::log.info "Done updating woob."
fi

bashio::log.info "Updating Woob dependencies..."
${KRUN} pip3 install $(${KRUN} python3 ./setup.py requirements)
bashio::log.info "Done updating Woob dependencies."

# ==============================================================================
# Set-up environment variables
# ==============================================================================

# Basic Kresus options
PORT="$(bashio::addon.ingress_port)"
HOST="$(bashio::addon.ip_address)"
KRESUS_URL_PREFIX="$(bashio::addon.ingress_entry)"
export PORT HOST KRESUS_URL_PREFIX

export KRESUS_PYTHON_EXEC=python3
export KRESUS_WOOB_DIR="${WOOB_DIR}"

# Kresus salt (generated once)
if ! bashio::fs.file_exists "/data/kresus_salt"; then
  pwgen 32 1 > /data/kresus_salt
  chmod 400 /data/kresus_salt
fi
KRESUS_SALT="$(cat /data/kresus_salt)"
export KRESUS_SALT

# Kresus Basic auth
if bashio::config.has_value 'http_basicauth'; then
  KRESUS_AUTH="$(bashio::config 'http_basicauth')"
  export KRESUS_AUTH
fi

# Kresus database
export KRESUS_DB_TYPE="postgres"

KRESUS_DB_HOST="$(bashio::config 'postgres_hostname')"
KRESUS_DB_PORT="$(bashio::config 'postgres_port')"
KRESUS_DB_USERNAME="$(bashio::config 'postgres_user')"
KRESUS_DB_PASSWORD="$(bashio::config 'postgres_password')"
KRESUS_DB_NAME="$(bashio::config 'postgres_database')"
export KRESUS_DB_HOST KRESUS_DB_PORT KRESUS_DB_USERNAME \
       KRESUS_DB_PASSWORD KRESUS_DB_NAME

# ==============================================================================
# Init data directory
# ==============================================================================
if ! bashio::fs.directory_exists "${KRESUS_DATA_DIR}"; then
  mkdir ${KRESUS_DATA_DIR}
  chown ${KRESUS_USER}:${KRESUS_USER} ${KRESUS_DATA_DIR}
  chmod 0700 ${KRESUS_DATA_DIR}
fi

# ==============================================================================
# Start Kresus
# ==============================================================================
${KRUN} kresus -c ${KRESUS_INI_FILE}
